# services/rss2email@.service
[Unit]
Description=rss2email caller for %i

[Service]
Type=oneshot
ExecStart=/opt/rss2email/scripts/renew.sh %i
